# home/.config/autostart/electron.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=electron
Comment=electronstartup script
Exec=/nix/store/rch0lw6p27vdhdayrjs60ldvf90gdgln-electron-unwrapped-40.9.2/libexec/electron/electron
StartupNotify=false
Terminal=false