# root/etc/systemd/system/lid-handler.service
[Unit]
Description=Handle lid close events

[Service]
Type=simple
ExecStart=~/Scripts/handle_lid_close.sh

[Install]
WantedBy=multi-user.target
